# deploy/cmd2/crontab/php_partying.cmd2.2.crontab
0 5 * * * php /home/ecs-user/webroot/php_partying/tutorial/app/cli.php main test cmd1 crontab2
