# deploy/cmd1/supervisor/php_partying.cl1.conf
[program:php_partying.cl1]
command=php /home/ecs-user/webroot/pt-oversea/tutorial/app/cli.php main test cmd1 universe; Command to run your program
autostart=false                         ; Start Supervisor on system boot
autorestart=false                       ; Restart program automatically if it crashes
stderr_logfile=/home/ecs-user/log/supervisor/test_cli.err   ; Error log file
stdout_logfile=/home/ecs-user/log/supervisor/test_cli.out   ; Output log file
